# supervisord.conf
[program:sugar]
command=/home/pi/envs/pi/bin/python /home/pi/apps/sugar/sugar.py
directory=/home/pi/apps/sugar
autostart=true
autorestart=true
redirect_stderr=true
